# root/etc/services.d/filebot/run
#!/usr/bin/with-contenv bash

FILEBOT_PORT=${FILEBOT_PORT:-7676}
FILEBOT_FORMAT=${FILEBOT_FORMAT:-'{plex}'}
FILEBOT_ACTION=${FILEBOT_ACTION:-duplicate}
FILEBOT_CONFLICT=${FILEBOT_CONFLICT:-auto}
FILEBOT_SERIES_DB=${FILEBOT_SERIES_DB:-TheTVDB}
FILEBOT_ANIME_DB=${FILEBOT_ANIME_DB:-AniDB}
FILEBOT_MOVIE_DB=${FILEBOT_MOVIE_DB:-TheMovieDB}
FILEBOT_MUSIC_DB=${FILEBOT_MUSIC_DB:-ID3}
SHARE_BASE_PATH=${SHARE_BASE_PATH:-/share}
UT_DIR='$v_path'
UT_TITLE='$v_name'
UT_LABEL='$v_label'

# necessary to take advantage of hardlinks
DIR_FULL_PATH="$SHARE_BASE_PATH$UT_DIR"
MEDIA_FULL_PATH="$SHARE_BASE_PATH/media"

exec \
  s6-setuidgid abc /app/shell2http \
    -port=${FILEBOT_PORT} \
    -show-errors \
    -include-stderr \
    -form /amc \
    "filebot -script fn:amc --output ${MEDIA_FULL_PATH} --action ${FILEBOT_ACTION} --conflict ${FILEBOT_CONFLICT} -no-xattr --log-file /config/amc.log \
    --def unsorted=n clean=y music=y artwork=y subtitles=en minLengthMS=0 \
    movieFormat=\"${FILEBOT_FORMAT}\" \
    seriesFormat=\"${FILEBOT_FORMAT}\" \
    animeFormat=\"${FILEBOT_FORMAT}\" \
    musicFormat=\"${FILEBOT_FORMAT}\" \
    seriesDB=\"${FILEBOT_SERIES_DB}\" \
    animeDB=\"${FILEBOT_ANIME_DB}\" \
    movieDB=\"${FILEBOT_MOVIE_DB}\" \
    musicDB=\"${FILEBOT_MUSIC_DB}\" \
    ut_dir=\"\$(eval echo '$DIR_FULL_PATH')\" \
    ut_kind=\"multi\" \
    ut_title=\"\$(eval echo '$UT_TITLE')\" \
    ut_label=\"\$(eval echo '$UT_LABEL')\""
